# host/portainer-update.service
[Unit]
Description=Portainer stack aktualisieren
After=docker.service network-online.target
Requires=docker.service

[Service]
Type=oneshot
{% set repo_dir = (playbook_dir ~ '/..') | realpath %}
ExecStartPre=/usr/bin/git -C {{ repo_dir }} -c safe.directory={{ repo_dir }} pull --ff-only
ExecStart=/usr/bin/docker compose -f {{ repo_dir }}/infrastructure/portainer/docker-compose.yml up -d --pull always
